# 99-reco-config-unmount.rules
ACTION=="remove", SUBSYSTEM=="block", KERNEL=="sd*", RUN+="/usr/bin/systemd-umount --no-block --collect  /media/RECO_CONFIG"
