# 99-reco-config.rules
ACTION=="add", SUBSYSTEM=="block", KERNEL=="sd*", RUN+="/usr/local/bin/mount_reco_config.sh %E{DEVNAME}"
