# systemd/reco_config_monitor.path
[Unit]
Description=Monitor /media/RECO_CONFIG for mount events

[Path]
PathExists=/media/RECO_CONFIG

[Install]
WantedBy=multi-user.target
